# .deploy/config/crontab
* * * * * /usr/local/bin/php /srv/app/crontab.php > /proc/1/fd/1 2>/proc/1/fd/2